# luci-app-ssrlocal/root/etc/uci-defaults/luci-ssrlocal
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ssrlocal[-1]
	add ucitrack ssrlocal
	set ucitrack.@ssrlocal[-1].init=ssrlocal
	commit ucitrack
EOF


rm -f /tmp/luci-indexcache
exit 0
